# please-review/from-root/ops/openclaw-tailscale-admin.service
[Unit]
Description=Ensure Tailscale SSH admin access for OpenClaw gateway
After=network-online.target tailscaled.service
Wants=network-online.target
ConditionPathExists=/home/marcu/ai-hell-mary/scripts/tailscale_admin_setup.sh

[Service]
Type=oneshot
ExecStart=/bin/bash /home/marcu/ai-hell-mary/scripts/tailscale_admin_setup.sh
RemainAfterExit=yes
TimeoutStartSec=90

[Install]
WantedBy=multi-user.target
